# scripts/com.ipark-drawing.notice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!--
    Wednesday 10:00 trigger. Polling logic lives in scripts/run-notice.sh
    (5 attempts at ~30 min intervals through ~12:00).
    launchd weekday: 0=Sun, 1=Mon, 2=Tue, 3=Wed, 4=Thu, 5=Fri, 6=Sat
  -->
  <key>Label</key>
  <string>com.ipark-drawing.notice</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/justinjeong/Workspace/ipark-drawing &amp;&amp; ./scripts/run-notice.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>3</integer>
    <key>Hour</key><integer>10</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/justinjeong/Workspace/ipark-drawing/data/notice.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/justinjeong/Workspace/ipark-drawing/data/notice.err.log</string>

  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
