# scripts/com.ipark-drawing.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!--
    Thursday-only morning run.
    Operator publishes the comment-target post on Thursday morning IF a vacancy
    exists; the bot auto-discovers today's article id and runs only when found.

    launchd weekday convention: 0=Sun, 1=Mon, 2=Tue, 3=Wed, 4=Thu, 5=Fri, 6=Sat
  -->
  <key>Label</key>
  <string>com.ipark-drawing.morning</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/justinjeong/Workspace/ipark-drawing &amp;&amp; ./scripts/run-morning.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>4</integer>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>57</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/justinjeong/Workspace/ipark-drawing/data/morning.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/justinjeong/Workspace/ipark-drawing/data/morning.err.log</string>

  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
